# docker/supervisor/conf.d/workers.conf
[program:queue-worker-default]
command=/bin/sh -lc "if php ./yii list --raw | grep -qx 'queue/listen'; then exec %(ENV_QUEUE_WORKER_1_COMMAND)s; else echo '[workers] queue/listen is not available, worker is idle'; exec tail -f /dev/null; fi"
process_name=%(program_name)s
autostart=true
autorestart=true
startsecs=3
numprocs=1
user=www-data
stdout_logfile=/app/runtime/logs/supervisor/queue-worker-default.stdout.log
stderr_logfile=/app/runtime/logs/supervisor/queue-worker-default.stderr.log
stopasgroup=true
killasgroup=true
group=workers

[program:queue-worker-priority]
command=/bin/sh -lc "if php ./yii list --raw | grep -qx 'queue/listen'; then exec %(ENV_QUEUE_WORKER_2_COMMAND)s; else echo '[workers] queue/listen is not available, worker is idle'; exec tail -f /dev/null; fi"
process_name=%(program_name)s
autostart=true
autorestart=true
startsecs=3
numprocs=1
user=www-data
stdout_logfile=/app/runtime/logs/supervisor/queue-worker-priority.stdout.log
stderr_logfile=/app/runtime/logs/supervisor/queue-worker-priority.stderr.log
stopasgroup=true
killasgroup=true
group=workers

[group:workers]
programs=queue-worker-default,queue-worker-priority
